# root/etc/services.d/terraria/run
#!/bin/sh
#cd /app/terraria/bin || exit
#s6-setuidgid terraria ./TerrariaServer.bin.x86_64 -config /config/serverconfig.txt
